# monitoring/node-exporter.service
[Unit]
Description=Node Exporter for Cybersecurity Training Platform
Documentation=https://prometheus.io/docs/guides/node-exporter/
After=network.target
Wants=network.target

[Service]
Type=simple
User=%i
Group=%i
ExecStart=/home/mark/sec-prac/monitoring/node_exporter \
    --web.listen-address="127.0.0.1:9100" \
    --path.rootfs=/host \
    --collector.filesystem.mount-points-exclude="^/(sys|proc|dev|host|etc)($$|/)" \
    --collector.textfile.directory=/var/lib/node_exporter/textfile_collector \
    --no-collector.infiniband \
    --no-collector.wifi \
    --no-collector.hwmon \
    --collector.systemd \
    --collector.processes
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=on-failure
RestartSec=5s
TimeoutStopSec=20s
SyslogIdentifier=node_exporter

# Security settings
NoNewPrivileges=yes
DynamicUser=yes
ProtectHome=yes
ProtectSystem=strict
ProtectControlGroups=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
PrivateDevices=yes
PrivateTmp=yes
ProtectClock=yes
ProtectHostname=yes
ProtectKernelLogs=yes
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Resource limits
LimitNOFILE=65536
LimitNPROC=32768

[Install]
WantedBy=default.target