# services/container-cleanup.service
[Unit]
Description=Container Resource Cleanup Service
Documentation=file:///home/mark/sec-prac/scripts/cleanup.sh
After=network.target
Wants=podman.service
BindsTo=user.slice

[Service]
Type=oneshot
ExecStart=/home/mark/sec-prac/scripts/cleanup.sh
User=%i
Group=%i

# Security settings
NoNewPrivileges=true
PrivateNetwork=false
PrivateDevices=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/mark/sec-prac
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true
RestrictNamespaces=~CLONE_NEWUSER

# Resource limits
MemoryLimit=256M
CPUQuota=25%
TasksMax=50

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=container-cleanup

# Environment
Environment="LOG_FILE=/home/mark/sec-prac/logs/container-cleanup.log"
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

[Install]
WantedBy=default.target

[Timer]
OnCalendar=*:0/15  # Run every 15 minutes
Persistent=true
RandomizedDelaySec=300  # Add up to 5 minutes random delay

[X-Container-Cleanup]
# Custom section for cleanup configuration
ContainerAgeHours=1
VolumeAgeHours=24
EnableSystemPrune=true
EnableTempCleanup=true