# services/container-cleanup.timer
[Unit]
Description=Container Resource Cleanup Timer
Documentation=file:///home/mark/sec-prac/scripts/cleanup.sh
Requires=container-cleanup.service

[Timer]
OnCalendar=*-*-* *:0/15  # Run every 15 minutes
Persistent=true
RandomizedDelaySec=300   # Add up to 5 minutes random delay to prevent resource spikes
AccuracySec=1min        # Allow 1 minute scheduling accuracy

# Cleanup scheduling
OnBootSec=10min         # Wait 10 minutes after boot
OnStartupSec=5min       # Wait 5 minutes after timer starts

[Install]
WantedBy=timers.target